# bb-radar.service
[Unit]
Description=BB Radar Service
After=network.target

[Service]
Type=simple
User=lain
WorkingDirectory=/home/lain/bb-radar
ExecStart=/home/lain/bb-radar/start.sh
Restart=on-failure
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

